# server/setup/api.service
[Unit]
Description=API Service
After=network.target

[Service]
ExecStart=/usr/local/bin/api
WorkingDirectory=/var/lib/api
Restart=always
User=root
Group=root
Environment=DB_PATH=/var/lib/api/items.db

[Install]
WantedBy=multi-user.target
